# Docker/supervisord.conf
[supervisord]
nodaemon=true

[program:atlas-engine]
command=atlas-engine

[program:bpmn-studio]
command=bpmn-studio --port=9000 --host=0.0.0.0
