# Docker/supervisord.conf
[supervisord]
nodaemon=true

[program:atlas-engine-runtime]
command=atlas-engine-runtime

[program:bpmn-studio]
command=bpmn-studio --port=9000 --host=0.0.0.0
